# .docker/app/base/drupal-cron.crontab
SHELL=/bin/bash
BASH_ENV=/container.env

0 * * * * runuser -u www-data -- /usr/local/bin/php /usr/app/vendor/bin/drush --root=/usr/app core:cron > /proc/1/fd/1 2>/proc/1/fd/2
# An empty line is required at the end of this file for a valid cron file.
